# services/systemd-units/memories-ttl-sweep.timer
[Unit]
Description=Hourly TTL sweep for memories.expires_at
Requires=memories-ttl-sweep.service

[Timer]
OnCalendar=*-*-* *:13:00 UTC
Persistent=true
RandomizedDelaySec=60
AccuracySec=30

[Install]
WantedBy=timers.target
